# SCRIPTS/weather.cron
0 * * * * /bin/bash -l -c 'ruby /Library/WebServer/Documents/RoR/TESTING/PayYourRent/SCRIPTS/get_weather_updates.rb 02110' >> /Library/WebServer/Documents/RoR/TESTING/PayYourRent/log/cron.log 2>&1 
0 * * * * /bin/bash -l -c 'ruby /Library/WebServer/Documents/RoR/TESTING/PayYourRent/SCRIPTS/get_weather_updates.rb 33871' >> /Library/WebServer/Documents/RoR/TESTING/PayYourRent/log/cron.log 2>&1 
0 * * * * /bin/bash -l -c 'ruby /Library/WebServer/Documents/RoR/TESTING/PayYourRent/SCRIPTS/get_weather_updates.rb 96801' >> /Library/WebServer/Documents/RoR/TESTING/PayYourRent/log/cron.log 2>&1 
0 * * * * /bin/bash -l -c 'ruby /Library/WebServer/Documents/RoR/TESTING/PayYourRent/SCRIPTS/get_weather_updates.rb 99524' >> /Library/WebServer/Documents/RoR/TESTING/PayYourRent/log/cron.log 2>&1 
0 * * * * /bin/bash -l -c 'ruby /Library/WebServer/Documents/RoR/TESTING/PayYourRent/SCRIPTS/get_weather_updates.rb 82941' >> /Library/WebServer/Documents/RoR/TESTING/PayYourRent/log/cron.log 2>&1 
0 * * * * /bin/bash -l -c 'ruby /Library/WebServer/Documents/RoR/TESTING/PayYourRent/SCRIPTS/get_weather_updates.rb 84321' >> /Library/WebServer/Documents/RoR/TESTING/PayYourRent/log/cron.log 2>&1 
0 * * * * /bin/bash -l -c 'ruby /Library/WebServer/Documents/RoR/TESTING/PayYourRent/SCRIPTS/get_weather_updates.rb 07039' >> /Library/WebServer/Documents/RoR/TESTING/PayYourRent/log/cron.log 2>&1 
0 * * * * /bin/bash -l -c 'ruby /Library/WebServer/Documents/RoR/TESTING/PayYourRent/SCRIPTS/get_weather_updates.rb 39530' >> /Library/WebServer/Documents/RoR/TESTING/PayYourRent/log/cron.log 2>&1 
0 * * * * /bin/bash -l -c 'ruby /Library/WebServer/Documents/RoR/TESTING/PayYourRent/SCRIPTS/get_weather_updates.rb 33710' >> /Library/WebServer/Documents/RoR/TESTING/PayYourRent/log/cron.log 2>&1 
0 * * * * /bin/bash -l -c 'ruby /Library/WebServer/Documents/RoR/TESTING/PayYourRent/SCRIPTS/get_weather_updates.rb 01844' >> /Library/WebServer/Documents/RoR/TESTING/PayYourRent/log/cron.log 2>&1 
